# bash/monitor.service
[Unit]
Description=Monitor File

[Service]
ExecStart=/bin/bash /usr/bin/monitor.sh


[Install]
WantedBy=multi-user.target