# bash/monitor.service
[Unit]
Description=Monitor File

[Service]
ExecStart=/usr/bin/monitor.sh

[Install]
WantedBy=multi-user.target